# oftcdns/cron.d/oftcdns.in
* 0 * * * root make -C @LOCALSTATEDIR@/lib/@PACKAGE@ @PACKAGE@.mak
