# oftcdns/cron.d/oftcdns.in
* 0 * * * root make -C @LOCALSTATEDIR@/@PACKAGE@ -f @PACKAGE@.mak
